# launchd/com.vibez-monitor.classify-missing.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.vibez-monitor.classify-missing</string>
    <key>ProgramArguments</key>
    <array>
        <string>__VIBEZ_ROOT__/backend/.venv/bin/python</string>
        <string>__VIBEZ_ROOT__/backend/scripts/classify_backfill.py</string>
        <string>--db-path</string>
        <string>__VIBEZ_ROOT__/vibez.db</string>
        <string>--model</string>
        <string>qwen2.5:3b</string>
        <string>--concurrency</string>
        <string>8</string>
        <string>--num-predict</string>
        <string>220</string>
        <string>--chunk-size</string>
        <string>100</string>
        <string>--lock-file</string>
        <string>/tmp/vibez-classify-backfill.lock</string>
    </array>
    <key>WorkingDirectory</key>
    <string>__VIBEZ_ROOT__</string>
    <key>RunAtLoad</key>
    <true/>
    <key>StartInterval</key>
    <integer>1800</integer>
    <key>StandardOutPath</key>
    <string>__LOG_DIR__/classify-missing-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>__LOG_DIR__/classify-missing-stderr.log</string>
</dict>
</plist>
