# launchd/com.vibez-monitor.synthesis.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.vibez-monitor.synthesis</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/braydon/projects/personal/vibez-monitor/backend/.venv/bin/python</string>
        <string>/Users/braydon/projects/personal/vibez-monitor/backend/scripts/run_synthesis.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/braydon/projects/personal/vibez-monitor</string>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/braydon/Library/Logs/vibez-monitor/synthesis-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/braydon/Library/Logs/vibez-monitor/synthesis-stderr.log</string>
</dict>
</plist>
