# launchd/com.vibez-monitor.sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.vibez-monitor.sync</string>
    <key>ProgramArguments</key>
    <array>
        <string>__VIBEZ_ROOT__/backend/.venv/bin/python</string>
        <string>__VIBEZ_ROOT__/backend/scripts/run_sync.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>__VIBEZ_ROOT__</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>VIBEZ_DB_PATH</key>
        <string>__VIBEZ_ROOT__/vibez.db</string>
        <key>BEEPER_API_URL</key>
        <string>http://localhost:23373</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>__LOG_DIR__/sync-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>__LOG_DIR__/sync-stderr.log</string>
</dict>
</plist>
