# launchd/com.vibez-monitor.dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.vibez-monitor.dashboard</string>

    <key>ProgramArguments</key>
    <array>
        <string>__VIBEZ_ROOT__/dashboard/scripts/start-dashboard.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__VIBEZ_ROOT__/dashboard</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__LOG_DIR__/dashboard-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>__LOG_DIR__/dashboard-stderr.log</string>
</dict>
</plist>
