# launchd/com.vibez-monitor.sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.vibez-monitor.sync</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/braydon/projects/personal/vibez-monitor/backend/.venv/bin/python</string>
        <string>/Users/braydon/projects/personal/vibez-monitor/backend/scripts/run_sync.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/braydon/projects/personal/vibez-monitor</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/braydon/Library/Logs/vibez-monitor/sync-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/braydon/Library/Logs/vibez-monitor/sync-stderr.log</string>
</dict>
</plist>
